# deploy/allinone/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=warn

[program:redis]
command=redis-server --maxmemory 128mb --maxmemory-policy allkeys-lru --save ""
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:backend]
command=uvicorn src.main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/app/backend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=ENVIRONMENT="production",ALLOWED_ORIGINS="http://localhost,http://localhost:80",REDIS_URL="redis://127.0.0.1:6379",NEC_WORKDIR="/tmp/nec_workdir",PYTHONUNBUFFERED="1"
priority=20

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
